# pkg/deb/prerm
#!/bin/sh

update-rc.d -f kmesg remove >/dev/null 2>&1 || : 

# uinstall
if [ -L /etc/init.d/kmesg ]; then
    /etc/init.d/kmesg stop >/dev/null 2>&1 || :
    rm -f /etc/init.d/kmesg || :
fi
